# src/roccat-vulcan.service
[Unit]
Description=Roccat Vulcan

[Service]
ExecStart=/usr/bin/roccat-vulcan

[Install]
WantedBy=default.target
